# .config/autostart/myclipboard.desktop
[Desktop Entry]
Type=Application
Name=Clipboard Sync
Exec=autocutsel -fork && autocutsel -selection PRIMARY -fork
Hidden=false
